# neuroinsight.service
[Unit]
Description=NeuroInsight MRI Processing System
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/src/desktop_alone_web_1
ExecStart=/home/ubuntu/src/desktop_alone_web_1/neuroinsight start
ExecStop=/home/ubuntu/src/desktop_alone_web_1/neuroinsight stop
Restart=always
RestartSec=10

# Environment variables
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=DOCKER_HOST=unix:///var/run/docker.sock

# Resource limits
MemoryLimit=4G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuroinsight

[Install]
WantedBy=multi-user.target